# systemd/dual-camera-record.service
[Unit]
Description=Dual 4K Camera QSV Recording Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=bjoern
Group=video
WorkingDirectory=/home/bjoern
Environment="LIBVA_DRIVER_NAME=iHD"
Environment="GST_VAAPI_ALL_DRIVERS=1"

# Wait for cameras to be ready
ExecStartPre=/bin/sleep 5
ExecStartPre=/bin/bash -c 'while [ ! -c /dev/video0 ] || [ ! -c /dev/video2 ]; do sleep 2; done'

# Main recording command
ExecStart=/home/bjoern/git/n100-videorecorder/scripts/dual-camera-record.sh

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dual-camera-record

# Security and permissions
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=false
RestrictRealtime=true

# Device access
DeviceAllow=/dev/dri/renderD128 rw
DeviceAllow=/dev/dri/card0 rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/video1 rw  
DeviceAllow=/dev/video2 rw
DeviceAllow=/dev/video3 rw
SupplementaryGroups=video render

# Resource limits (adjust based on your system)
LimitNOFILE=4096
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
